# src/bin/gunicorn/gunicorn.service
[Unit]
Description=gunicorn daemon
Requires=gunicorn.socket
After=network.target

[Service]
PIDFile=/run/gunicorn/pid
User=zod
Group=users
RuntimeDirectory=gunicorn
WorkingDirectory=/var/www/castic/src
ExecStart=/var/www/castic/src/bin/gunicorn/gunicorn.sh	
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s TERM $MAINPID
PrivateTmp=true

[Install]
WantedBy=multi-user.target
